# net/keepalived/files/etc/hotplug.d/keepalived/700-luci
#!/bin/sh

. /lib/functions/keepalived/hotplug.sh

add_sync_file /etc/config/luci

keepalived_hotplug
