# net/libreswan/files/etc/hotplug.d/libreswan/02-vti
#!/bin/sh

addvti() {
	[ -d /proc/sys/net/ipv4/conf/${VTI_IFACE} ] && ip tun del ${VTI_IFACE}

	ip tunnel add ${VTI_IFACE} mode vti local ${PLUTO_ME} remote ${PLUTO_PEER} okey ${CONNMARK_OUT%/*} ikey ${CONNMARK_IN%/*}  

	sysctl -w net.ipv4.conf.${VTI_IFACE}.disable_policy=1
	sysctl -w net.ipv4.conf.${VTI_IFACE}.rp_filter=0
	sysctl -w net.ipv4.conf.${VTI_IFACE}.forwarding=1

	[ -n "${VTI_IP}" ] && ip addr add ${VTI_IP} dev ${VTI_IFACE}

	ip link set ${VTI_IFACE} up
}

delvti() {
	[ -n "${VTI_IFACE}" -a -d /proc/sys/net/ipv4/conf/${VTI_IFACE} ] && ip tun del ${VTI_IFACE}
}

[ -z "$VTI_IFACE" ] && exit 0

case "${PLUTO_VERB}" in
	up-host|up-client) addvti ;;
	down-host|down-client) delvti ;;
esac
